# dot_config/systemd/user/update-weather.service
[Unit]
Description=Script to update weather for i3bar

[Service]
Type=simple
ExecStart=%h/bin/update_weather

[Install]
WantedBy=default.target
