# dot_config/systemd/user/update-solar.service
[Unit]
Description=Script to update solar panels for i3bar

[Service]
Type=simple
ExecStart=%h/bin/update_egauge

[Install]
WantedBy=default.target
